# .config/autostart/qbittorrent.desktop
[Desktop Entry]
Type=Application
Exec=flatpak run org.qbittorrent.qBittorrent
Name=qBittorrent
