# .config/autostart/redshift.desktop
[Desktop Entry]
Type=Application
Exec=redshift
Name=Redshift
