# .config/autostart/redshift.desktop
[Desktop Entry]
Type=Application
Exec=redshift -l 47.4385:18.9108
Name=Redshift
