# .config/autostart/webcord.desktop
[Desktop Entry]
Type=Application
Exec=flatpak run io.github.spacingbat3.webcord
Name=WebCord
